# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/ghost-tunnel/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── sshd ─────────────────────────────────────────
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
startretries=5
stopwaitsecs=5
stdout_logfile=/var/log/ghost-tunnel/sshd.log
stderr_logfile=/var/log/ghost-tunnel/sshd.err
stdout_logfile_maxbytes=5MB
priority=10

; ── health check HTTP server ──────────────────────
[program:health]
command=python3 /usr/local/bin/health.py
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/ghost-tunnel/health.log
stderr_logfile=/var/log/ghost-tunnel/health.err
stdout_logfile_maxbytes=2MB
priority=20

; ── bore tunnel manager ───────────────────────────
; startretries=999 = effectively unlimited reconnect
[program:tunnel]
command=/usr/local/bin/tunnel.sh
autostart=true
autorestart=true
startretries=999
stopwaitsecs=10
stdout_logfile=/var/log/ghost-tunnel/tunnel.log
stderr_logfile=/var/log/ghost-tunnel/tunnel.err
stdout_logfile_maxbytes=10MB
priority=30

; ── system watchdog ───────────────────────────────
[program:watchdog]
command=/usr/local/bin/watchdog.sh
autostart=true
autorestart=true
startretries=999
stdout_logfile=/var/log/ghost-tunnel/watchdog.log
stderr_logfile=/var/log/ghost-tunnel/watchdog.err
stdout_logfile_maxbytes=5MB
priority=99
